# dist/plugin_loader-prerelease.service
[Unit]
Description=SteamDeck Plugin Loader
After=network-online.target
Wants=network-online.target
[Service]
Type=simple
User=root
Restart=always
ExecStart=${HOME}/homebrew/services/PluginLoader
WorkingDirectory=${HOME}/homebrew/services
Environment=PLUGIN_PATH=${HOME}/homebrew/plugins
Environment=LOG_LEVEL=DEBUG
[Install]
WantedBy=multi-user.target